# systemd/raspi-hotspot.service
[Unit]
Description=Raspi Hotspot — WiFi Access Point (uap0)
After=network.target
Before=raspi-app.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/home/pi/raspi-hotspot/config/.env
ExecStart=/usr/local/bin/raspi-hotspot-up
ExecStop=/bin/bash -c "\
    pkill -f 'hostapd /etc/hostapd/hostapd.conf' || true; \
    pkill -f 'raspi-dnsmasq.pid' || true; \
    iw dev uap0 del 2>/dev/null || true"

[Install]
WantedBy=multi-user.target
