# systemd/raspi-hotspot.service
[Unit]
Description=Raspi Hotspot — WiFi Access Point
After=network.target
Before=raspi-app.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/home/pi/raspi-hotspot/config/.env
ExecStart=/usr/local/bin/raspi-hotspot-up
ExecStop=/usr/local/bin/raspi-hotspot-down

[Install]
WantedBy=multi-user.target
